# infra/shorts-pipeline.service
[Unit]
Description=Shorts Bot Automated Pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/home/ec2-user/shorts-bot

ExecStart=/usr/bin/timeout 5400 \
    /home/ec2-user/shorts-bot/.venv/bin/python -m src.main run-pipeline

ExecStopPost=/bin/bash -c \
    'REGION=$(curl -sf http://169.254.169.254/latest/meta-data/placement/region || echo us-east-1); \
     aws ec2 stop-instances \
     --instance-ids "$(curl -sf http://169.254.169.254/latest/meta-data/instance-id)" \
     --region "$REGION" || true'

EnvironmentFile=/home/ec2-user/shorts-bot/.env
StandardOutput=journal
StandardError=journal
Restart=no

[Install]
WantedBy=multi-user.target
